# setup/flavors/rkt/mailu-resolver.service
{% set env='mailu.env' %}
[Unit]
# This file is auto-generated by the Mailu configuration wizard.
# Please read the documentation before attempting any change.
# Generated for {{ flavor }} flavor
Description=Mailu resolver component

[Service]

Restart=always

EnvironmentFile={{ root }}/{{ env }}

ExecStart=/usr/bin/rkt --inherit-env --dns 8.8.8.8 --uuid-file-save="{{ root }}/resolver.uuid" \
    mailu/unbound:{{ version }}
    
ExecStopPost=/bin/sh -c "rkt rm $(cat {{ root }}/resolver.uuid)"
