# setup/flavors/rkt/smtp.service
{% set env='mailu.env' %}
[Unit]
# This file is auto-generated by the Mailu configuration wizard.
# Please read the documentation before attempting any change.
# Generated for {{ flavor }} flavor
Description=Mailu smtp component
Requires=mailu-front.service
{% if resolver_enabled %}
Requires=mailu-resolver.service

[Service]

Restart=always

EnvironmentFile={{ root }}/{{ env }}

ExecStart=/usr/bin/rkt --inherit-env --dns {{ dns }} --uuid-file-save="{{ root }}/smtp.uuid" \
    --volume overrides,kind=host={{ root }}/overrides \
    ${DOCKER_ORG:-mailu}/postfix:${MAILU_VERSION:-{{ version }}} \
    --mount volume=overrides,target=/overrides
    
ExecStopPost=/bin/sh -c "rkt rm $(cat {{ root }}/smtp.uuid)"
